# dist/linux/tagent_init.service
[Unit]
Description=TrustAgent Initialization

[Service]
Type=oneshot
WorkingDirectory=/usr/bin/trustagent
ExecStart=/usr/bin/trustagent/tagent init

[Install]
WantedBy=multi-user.target
